# supervisord.conf
; ============================================================
; supervisord.conf — CachyOS Headless service orchestration
; ============================================================
;
; Key design decisions:
;   - NO Xvfb:  TigerVNC is itself an X server, so a separate
;               virtual framebuffer is redundant and would clash
;               on the same DISPLAY.
;   - NO x11vnc: TigerVNC already serves VNC on port 5901.
;                Adding x11vnc would create a port conflict.
;   - PipeWire + pipewire-pulse: both are started so that
;                PulseAudio-compatible apps (games, browsers,
;                Discord) continue to work out of the box.
;   - pipewire-pulse uses LD_PRELOAD wrapper so it runs as the
;                same process tree as pipewire under the user.
; ============================================================

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

; ---- Audio ------------------------------------------------

[program:pipewire]
command=pipewire
autorestart=true
priority=10
user=cachyos
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:pipewire-pulse]
command=pipewire -c pipewire-pulse.conf
autorestart=true
priority=11
user=cachyos
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:wireplumber]
command=wireplumber
autorestart=true
priority=12
user=cachyos
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ---- VNC / Desktop ----------------------------------------

[program:vnc]
; TigerVNC creates its own X server on :1 — no Xvfb needed.
; Resolution and depth are read from the VNC_GEOMETRY and
; VNC_DEPTH env vars set in start.sh.
command=su - cachyos -c "vncserver :1 -geometry %(ENV_VNC_GEOMETRY)s -depth %(ENV_VNC_DEPTH)s -localhost no -rfbauth /home/cachyos/.vnc/passwd"
autorestart=true
priority=20
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:xfce]
command=su - cachyos -c "DISPLAY=:1 startxfce4"
autorestart=true
priority=30
environment=DISPLAY=":1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ---- noVNC / Web Access -----------------------------------

[program:novnc]
; websockify installed in /opt/noVNC-env virtual environment
command=/opt/noVNC-env/bin/websockify --web=/usr/share/novnc 8080 localhost:5901
autorestart=true
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
