# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:xvfb]
command=Xvfb :1 -screen 0 1920x1080x24 +extension GLX +render -noreset
autorestart=true
priority=10

[program:vnc]
command=su - cachyos -c "vncserver :1 -geometry 1920x1080 -depth 24 -rfbauth /home/cachyos/.vnc/passwd -rfbport 5901 -noxstartup"
autorestart=true
priority=20
startsecs=5

[program:x11vnc]
command=x11vnc -display :1 -forever -shared -nopw -rfbport 5901
autorestart=true
priority=30
startsecs=5

[program:xfce]
command=su - cachyos -c "DISPLAY=:1 startxfce4"
autorestart=true
priority=40
environment=DISPLAY=":1"

[program:pulseaudio]
command=pulseaudio --system --daemonize=no --log-target=stderr
autorestart=true
priority=15

[program:novnc]
command=websockify --web=/usr/share/novnc --cert=none --ssl-only=false 8080 localhost:5901
autorestart=true
priority=50
