# etc/supervisor.d/novnc.ini
[program:novnc]
priority=30
autostart=true
autorestart=true
user=root
command=/opt/noVNC-env/bin/websockify --web=/usr/share/novnc %(ENV_PORT_NOVNC_WEB)s localhost:%(ENV_PORT_VNC)s
stopsignal=INT
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
