# etc/supervisor.d/pipewire.ini
[program:pipewire]
priority=40
autostart=true
autorestart=true
user=%(ENV_USER)s
directory=/home/%(ENV_USER)s
command=/usr/bin/start-pipewire.sh
environment=HOME="/home/%(ENV_USER)s",USER="%(ENV_USER)s",DISPLAY="%(ENV_DISPLAY)s",PUID="%(ENV_PUID)s"
stopsignal=INT
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:wireplumber]
priority=41
autostart=true
autorestart=true
user=%(ENV_USER)s
directory=/home/%(ENV_USER)s
command=/usr/bin/start-wireplumber.sh
environment=HOME="/home/%(ENV_USER)s",USER="%(ENV_USER)s",DISPLAY="%(ENV_DISPLAY)s",PUID="%(ENV_PUID)s"
stopsignal=INT
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:pipewire-pulse]
priority=42
autostart=true
autorestart=true
user=%(ENV_USER)s
directory=/home/%(ENV_USER)s
command=/usr/bin/start-pipewire-pulse.sh
environment=HOME="/home/%(ENV_USER)s",USER="%(ENV_USER)s",DISPLAY="%(ENV_DISPLAY)s",PUID="%(ENV_PUID)s"
stopsignal=INT
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
