# etc/supervisor.d/x11vnc.ini
[program:x11vnc]
priority=30
autostart=true
autorestart=true
user=root
command=/usr/bin/start-x11vnc.sh
environment=DISPLAY="%(ENV_DISPLAY)s",HOME="/home/%(ENV_USER)s",PORT_VNC="%(ENV_PORT_VNC)s"
stopsignal=INT
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
